# etc/supervisord.d/grafana.ini
[program:grafana]
directory=/usr/share/grafana/
command = /usr/sbin/grafana-server --pidfile=/var/run/grafana-server.pid --config=/etc/grafana/grafana.ini cfg:default.paths.data=/var/lib/grafana cfg:default.paths.logs=/var/log/grafana
redirect_stderr=true
stdout_syslog=true
